# srcs/telegraf/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:telegraf]
command= exec /usr/bin/telegraf
autostart=true
autorestart=true